# tmp/droid-config.files
%defattr(-,root,root,-)
/etc/acdbdata
/etc/audio_platform_info.xml
/etc/gps.conf
/etc/izat.conf
/etc/lowi.conf
/etc/media_codecs.xml
/etc/media_codecs_ffmpeg.xml
/etc/media_codecs_google_audio.xml
/etc/media_codecs_google_telephony.xml
/etc/media_codecs_google_video.xml
/etc/media_codecs_performance.xml
/etc/media_profiles.xml
/etc/mixer_paths.xml
/etc/modprobe.d/blacklist-evbug.conf
/etc/obexd/excludes/bluetooth:pcsuite
/etc/obexd/noplugins/ftp
/etc/obexd/noplugins/irmc
/etc/obexd/noplugins/mas
/etc/obexd/noplugins/syncevolution
/etc/obexd/opp.conf
/etc/obexd/pbap.conf
/etc/obexd/root
/etc/ofono/binder.conf
/etc/ofono/binder.d/dual-sim.conf
/etc/sap.conf
/etc/sec_config
/etc/selinux/config
/etc/sensorfw/primaryuse.conf
/etc/sysctl.d/ipv4-settings.conf
/etc/sysctl.d/ipv6-settings.conf
/etc/usb-moded/dyn-modes/mtp_mode-droid.ini
/etc/xdg/QtProject/QPlatformTheme.conf
/etc/xdg/QtProject/QtQuick2.conf
/etc/xdg/QtProject/Sensors.conf
/etc/xtra_root_cert.pem
/etc/xtwifi.conf
/lib/udev/rules.d/997-droid-drm.rules
/lib/udev/rules.d/998-droid-system.rules
/usr/bin/droid/droid-hal-shutdown.sh
/usr/bin/droid/droid-hal-startup.sh
/usr/bin/droid/droid-init-done.sh
/usr/bin/droid/droid-load-firmware.sh
/usr/bin/droid/get_partition.sh
/usr/lib/systemd/system/actdead.target.wants/droid-hal-init.service
/usr/lib/systemd/system/droid-hal-init.service
/usr/lib/systemd/system/droid-late-start.service
/usr/lib/systemd/system/graphical.target.wants/droid-hal-init.service
/usr/lib/systemd/system/graphical.target.wants/droid-late-start.service
/usr/lib/tmpfiles.d/droid-hal.conf
/var/lib/environment/compositor/droid-hal-device.conf
/var/lib/environment/nemo/60-multimedia.conf
/var/lib/environment/ofono/noplugin.conf
/boot/extracting-README.txt
/boot/flashing-README.txt
/etc/acdbdata
/etc/gbinder.conf
/etc/gst-droid/gstdroidcodec.conf
/etc/mce/60-fingerprint-sensor.conf
/etc/mce/60-hybris-features.ini
/etc/mce/60-proximity-sensor.conf
/etc/ofono/ril_subscription.conf
/etc/ofono/ril_subscription.d/dual-sim.conf
/etc/sailfish-fpd/50-settings.ini
/usr/bin/droid/droid-bootctl.sh
/usr/bin/droid/droid-hal-early-init.sh
/usr/lib/systemd/system/droid-bootctl.service
/usr/lib/systemd/system/droid-hal-init.service
/usr/lib/systemd/system/graphical.target.wants/droid-bootctl.service
/usr/lib/systemd/system/multi-user.target.wants/wait_for_keymaster.service
/usr/lib/systemd/system/sailfish-unlock-agent.service.d/50-vendor.hwcomposer-2-3.conf
/usr/lib/systemd/system/sailfish-upgrade-ui.service.d/50-vendor.hwcomposer-2-3.conf
/usr/lib/systemd/system/wait_for_keymaster.service
/usr/lib/systemd/user/jolla-startupwizard-pre-user-session.service.d/50-vendor.hwcomposer-2-3.conf
/usr/lib/systemd/user/lipstick.service.d/50-vendor.hwcomposer-2-3.conf
/usr/libexec/droid-hybris/system/etc/init/disabled_services.rc
/usr/libexec/droid-hybris/system/etc/init/ecclist.rc
/usr/libexec/droid-hybris/system/etc/ld.config.29.txt
/usr/share/ngfd/plugins.d/60-profile.ini
/var/lib/environment/nemo/99-qtscenegraph.conf
/var/lib/environment/compositor/droid-hal-device.conf
/usr/bin/droid/droid-hal-early-init.sh
/bin/parse-android-dynparts
/usr/share/ssu/features.d/adaptation-community.ini
